# conditional-shutdown.service
[Unit]
Description=Shutdown the Minecraft server if there has been no recent activity.

[Service]
Type=oneshot
ExecStart=/StopServerWithNoPlayers
